# dot_config/systemd/user/claude-notify-cleanup.service
[Unit]
Description=Cleanup Claude Code notify dispatcher state files (F-3.next #5)
Documentation=https://github.com/kkiyama117/dotfiles
ConditionPathExists=%t/claude-notify/sessions

[Service]
Type=oneshot
ExecStart=%h/.local/bin/claude-notify-cleanup
# Defence in depth: read-only filesystem except for the runtime state dir
# the helper actually prunes.
ProtectSystem=strict
ProtectHome=read-only
ReadWritePaths=%t/claude-notify
PrivateTmp=true
NoNewPrivileges=true
